# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:postgresql]
command=/bin/bash -c "PGDATA=/data/postgresql /usr/bin/postgres -D /data/postgresql"
user=postgres
autostart=true
autorestart=true
stdout_logfile=/data/logs/postgresql-stdout.log
stderr_logfile=/data/logs/postgresql-stderr.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
priority=10
startsecs=10
startretries=5
environment=PGDATA="/data/postgresql",PGPORT="5432"

[program:go-backend]
command=/app/nginxops
directory=/app
autostart=true
autorestart=true
stdout_logfile=/data/logs/app/stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_syslog=true
redirect_stderr=true
priority=20
startsecs=10
startretries=5
environment=DATA_DIR=/data

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/data/logs/nginx/supervisor-stdout.log
stderr_logfile=/data/logs/nginx/supervisor-stderr.log
priority=30
startsecs=5
startretries=3
environment=DATA_DIR=/data
